# supervisord.conf
[program:l2tp]
command=/startup.sh
[program:client]
command=/client.sh
